# cron.d/postgres.cron
# m h  dom mon dow   command

MAILTO=fashiolista@9apps.net

date=$(date +"%Y%m%d-%H%M%S")
@daily /root/pgRDS/backup.sh basebackup ${date} > /dev/null 2>&1
